# deploy/systemd/reify-warm-lane.service
[Unit]
Description=Reify warm-lane CoW pool: provision+mount the XFS-reflink loopback volume (boot-persistent, fail-open)
Documentation=file:///home/leo/src/reify/docs/prds/warm-lane-pool-activation-seam.md
# NO After=local-fs.target: that is a system-manager target, not present in the
# per-user systemd manager.  User units start after the system's basic.target
# automatically; no explicit ordering directive is needed here.
# NO Requires= here — the orchestrator drop-in owns the soft dependency
# (Wants=reify-warm-lane.service in orchestrator-reify.service.d/warm-lane.conf).

[Service]
Type=oneshot
RemainAfterExit=yes
# Reuse the full idempotent #4659 provision path (P1 never-reformat + P2 fail-closed
# reflink probe).  Unit-active ⟺ mounted + reflink probe passed.
# ExecStart uses the stable main-checkout absolute path — host convention matching
# the live orchestrator unit and jcodemunch-serve.service.  Invoked with defaults:
# /media/leo/data_lv_1/leo/reify-warm-lanes.img, 4096 GiB (final sizing/location deferred to R4/#4665
# per PRD §9 Q3).
#
# PREREQUISITE — sudo NOPASSWD: provision-warm-lane-fs.sh routes mkfs.xfs/losetup/
# mount/chown through `sudo` when EUID != 0.  In a non-interactive user-service boot
# context there is no TTY, so sudo must be configured NOPASSWD for those commands:
#   /etc/sudoers.d/<user>: <user> ALL=(root) NOPASSWD: /sbin/losetup, /sbin/mkfs.xfs, \
#     /bin/mount, /bin/umount, /bin/chown
# Without NOPASSWD the oneshot fails; the Wants= (not Requires=) fail-open design
# ensures this degrades to the cold path (inv.6) without blocking orchestrator startup.
ExecStart=/home/leo/src/reify/scripts/provision-warm-lane-fs.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
