# deploy/systemd/orchestrator-reify.service.d/warm-lane.conf
# drop-in: orchestrator-reify.service.d/warm-lane.conf
#
# DA5 fail-open: Wants= (NOT Requires=) so a missing/failed warm-lane mount
# degrades to the cold path (inv.6) and never blocks orchestrator startup.
# After= establishes ordering (orchestrator starts after provisioning completes)
# without making the mount a hard dependency.
#
# This drop-in is owned by the reify repo and installed by
# scripts/install-warm-lane-units.sh.  It adds the warm-lane service to the
# orchestrator's After/Wants additively, without editing the dark-factory-owned
# base orchestrator-reify.service unit.

[Unit]
Wants=reify-warm-lane.service
After=reify-warm-lane.service
